# deploy/telegram2qbittorrent.service
[Unit]
Description=Telegram to qBittorrent bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__REPO__
ExecStart=__REPO__/venv/bin/python -u main.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
